# usr/lib/systemd/user/miatoll_data_quirk.service
[Unit]
Description=Miatoll data quirk

[Service]
ExecStart=/usr/bin/miatoll_data_quirk.sh

[Install]
WantedBy=multi-user.target
