# usr/lib/systemd/user/miatoll_quirks.service
[Unit]
Description=Miatoll quirks

[Service]
ExecStart=/usr/bin/sleep 10
ExecStart=/usr/bin/miatoll_quirks.sh

[Install]
WantedBy=default.target
