# usr/lib/systemd/user/mobile-quirks.service
[Unit]
Description=Mobile quirks

[Service]
ExecStartPre=/usr/bin/sleep 10
ExecStart=/usr/bin/mobile-quirks.sh

[Install]
WantedBy=default.target
