# usr/lib/systemd/system/miatoll-data-quirk.timer
[Unit]
Description=Miatoll data quirk

[Timer]
OnActiveSec=10
OnUnitActiveSec=120
Unit=miatoll-data-quirk.service

[Install]
WantedBy=timers.target
