# usr/lib/systemd/system/miatoll-data-quirk.service
[Unit]
Description=Miatoll data quirk

[Service]
ExecStart=/usr/bin/miatoll-data-quirk.sh
