# usr/lib/systemd/system/mobile-data-quirk.service
[Unit]
Description=Mobile data quirk

[Service]
ExecStart=/usr/bin/mobile-data-quirk.sh
