# usr/lib/systemd/system/ModemManager.service.d/11-miatoll.conf
[Service]
ExecStartPre=/usr/bin/miatoll-enable-modem.sh
