# usr/lib/systemd/system/ModemManager.service.d/11-ofono.conf
[Service]
ExecStartPre=/usr/bin/gdbus wait -y org.ofono
ExecStartPre=/usr/bin/sleep 2
