# etc/systemd/system/ModemManager.service.d/10-ofono.conf
[Unit]
After=ofono.service
Requires=ofono.service
